# systemd/user/logtime.service
[Unit]
Description=Persistent Go logger that prints the time every minute
After=default.target

[Service]
ExecStart=%h/.local/bin/logtime
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target

