# systemd/user/monitors.service

[Unit]
Description=Persistent monitor layout manager
After=graphical-session.target

[Service]
ExecStart=/usr/local/bin/monitors
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target

