# systemd/user/monitors.service
[Unit]
Description=Persistent monitor layout manager
After=graphical-session.target

[Service]
ExecStart=%h/.config/hypr/scripts/monitors.sh
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target

